# backup/crontab
0 2 * * 1 /bin/sh /backup.sh >> /var/log/cron.log 2>&1
